# config/supervisor/apache2.conf
# Supervisor configuration for Apache
[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autorestart=true
startretries=3
stdout_logfile=/var/log/apache2/supervisor.log
redirect_stderr=true